# supervisord.conf
[supervisord]
nodaemon=true
user=root
# Directing main supervisor logs to stdout so you see startup errors in HF Logs
logfile=/dev/stdout
logfile_maxbytes=0
pidfile=/tmp/supervisord.pid

[program:php-fpm]
command=php-fpm
autostart=true
autorestart=true

[program:nginx]
command=nginx -g 'daemon off;'
autostart=true
autorestart=true

[program:bot]
command=/app/bot/banabot
directory=/app/bot
autostart=true
autorestart=true
# Essential for scanning the WhatsApp QR code
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
